# systemd/xdg-desktop-portal-gtk.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/xdg-desktop-portal-gtk -r

[Install]
WantedBy=sway-session.target
